# docker/development/workspace/supervisord.conf
[supervisord]
nodaemon=true

[program:artisan-tinker]
command=php /var/www/artisan tinker
user=www
autostart=false
autorestart=false

# Ajoute ici d'autres commandes artisan/npm si besoin
